# thirdparty/Supervisor/auth-celerybeat.conf
[program:auth-celerybeat]
command=python manage.py celerybeat
directory=/home/allianceserver/allianceauth
user=allianceserver
stdout_logfile=/home/allianceserver/allianceauth/log/beat.log
stderr_logfile=/home/allianceserver/allianceauth/log/beat.log
autostart=true
autorestart=true
startsecs=10
priority=999
