# thirdparty/Supervisor/auth-celeryd.conf
[program:auth-celeryd]
command=python manage.py celeryd
directory=/home/allianceserver/allianceauth
user=nobody
numprocs=1
stdout_logfile=/home/allianceserver/allianceauth/log/worker.log
stderr_logfile=/home/allianceserver/allianceauth/log/worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=1000
